# eip-platform.service
[Unit]
Description=EIP Platform — Universal App Server
After=network.target

[Service]
User=nojus
WorkingDirectory=/home/nojus/eip-platform
EnvironmentFile=/home/nojus/eip-platform/.env
ExecStart=/home/nojus/eip-platform/.venv/bin/python /home/nojus/eip-platform/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
